# Resources/com.clawmail.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clawmail.agent</string>
    <key>Program</key>
    <string>/Applications/ClawMail.app/Contents/MacOS/ClawMailApp</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/ClawMail.app/Contents/MacOS/ClawMailApp</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/clawmail.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/clawmail.stderr.log</string>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
